# kde-steamdeck/etc/xdg/autostart/ibus.desktop
[Desktop Entry]
Name=IBus
Comment=Input method framework daemon
OnlyShowIn=KDE
Type=Application
Exec=ibus-daemon --desktop=kde
